# Linux/Systemd/timers/certbot-renew/certbot-renew.timer
[Unit]
Description=Renew TLS certificates using certbot renew

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
